# basefiles/tappd-socket.socket
# SPDX-FileCopyrightText: 2025 Phala Network <dstack@phala.network>
#
# SPDX-License-Identifier: Apache-2.0

# Backward compatibility socket for containers that mount /var/run/tappd.sock directly.
# The socket is owned by systemd, so it survives service restarts without inode changes.

[Unit]
Description=dstack backward compatibility socket (tappd.sock)

[Socket]
ListenStream=/var/run/tappd.sock
SocketMode=0777

[Install]
WantedBy=sockets.target
